# systemd/user/cursor-update.service
[Unit]
Description=Update Cursor AppImage
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
Environment=DOTFILES_DIR=%h/dotfiles
ExecStart=/bin/sh -c '. "$DOTFILES_DIR/core/managers/updateman/core/updateman.sh"; updateman cursor run'
